# .jobs/crontab
30 4 * * * ~/.shortcuts/tasks/syncthing.sh
