# .jobs/crontab
30 4 * * * ~/.shortcuts/tasks/syncthing.sh
0  5 * * * fish -c refreshLibrary
30 9 * * * ssh xk@192.168.1.114 lt-wakeup
