# .jobs/crontab
0 * * * * ~/.shortcuts/tasks/syncthing.sh

